# system_config/getty@tty1.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin alex --noclear %I 38400 linux

[Install]
WantedBy=multi-user.target


